# misc/services/kdeconnectd.service
[Unit]
Description=KDE Connect Daemon

[Service]
ExecStart=/usr/lib/kdeconnectd
Restart=on-failure

[Install]
WantedBy=default.target
