# misc/services/unclutter.service
[Unit]
Description=Unclutter mouse cursor when idle

[Service]
ExecStart=/usr/bin/unclutter
Restart=on-failure

[Install]
WantedBy=default.target
