# home/.config/systemd/user/cloudflared.service
[Unit]
Description=Cloudflare Tunnel
After=network.target

[Service]
TimeoutStartSec=0
Type=notify
ExecStart=/usr/bin/cloudflared tunnel run kuroneko
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
